# files/system/usr/lib/systemd/user/origami-updates.service
[Unit]
Description=Run Topgrade System Updates

[Service]
Type=oneshot
ExecStart=/usr/bin/topgrade -y --no-retry --disable firmware --skip-notify
